# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#

nohup python /home/debian/bbgdemo/iot_demo.py >/dev/null 2>&1 &
sleep 5
cd /home/debian/bbgdemo
nohup java -jar bbg-demo-1.0.jar -level severe >/dev/null 2>&1 &

exit 0
