# java/etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#

cd /home/debian/bbgdemo
sudo ./bbg_python_service.sh start
sleep 2
sudo ./bbg_java_service.sh start

exit 0